# 10-md407.rules
ATTR{product}=="FT230X Basic UART", RUN+="/usr/bin/set_md407_baud %k", SYMLINK+="md407c%1"